# docker/etc/services.d/20-celery-worker/finish
#!/usr/bin/with-contenv bash

cd /app
s6-setuidgid abc celery -A plaxt control shutdown
